# conf/tornado.ini
[program:my_tornado]
command=/home/codeboy/envs/main/bin/python /home/codeboy/projects/tornado_main/tornado_main.py